# linux/TP_03/systemd/units/backup.timer
[Unit]
Description=Run backup every hour

[Timer]
OnUnitActiveSec=59min


[Install]
WantedBy=timers.target